# gensyn.service
[Unit]
Description=Gensyn RL Swarm Node
After=network.target docker.service
Wants=docker.service

[Service]
User=root
WorkingDirectory=/home/gensyn/rl_swarm
ExecStart=/bin/bash /home/gensyn/rl_swarm/run_node.sh
Restart=always
RestartSec=10
LimitNOFILE=65535
# Optional: increase FileDescriptorLimit / MemoryLimit etc.

[Install]
WantedBy=multi-user.target
